# docker/supervisor.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /home/docker/code/docker/uwsgi.ini

[program:app]
command = /usr/sbin/nginx
